# server/python-ci.service.in
[Unit]
Description=Python-CI server
After=network.target

[Service]
Type=simple
User=?USER?
WorkingDirectory=?DIR?
ExecStart=?DIR?/start.sh

[Install]
WantedBy=multi-user.target
